# systemd/board-game-backup.timer
[Unit]
Description=Daily Board Game database backup

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
